# conf/rc.local.sample
# sample rc.local 
# add your code below
cp -a /etc/motd /var/run/motd
